# control_robot.service
[Unit]
Description=Control Robot Service
After=network.target docker.service
Wants=docker.service
Requires=network.target

[Service]
Type=simple
User=unitree
Group=unitree
WorkingDirectory=/home/unitree/control_robot
ExecStart=/home/unitree/control_robot/start_h1_unified.sh
ExecStop=/home/unitree/control_robot/stop_h1.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=control_robot

# Переменные окружения
Environment=NODE_ENV=production
Environment=HOME=/home/unitree
Environment=PATH=/home/unitree/.nvm/versions/node/v18.19.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NVM_DIR=/home/unitree/.nvm
Environment=VIRTUAL_ENV=/home/unitree/control_robot/backend/src/services/.venv

# Права доступа
ReadWritePaths=/home/unitree
ReadWritePaths=/home/unitree/control_robot
ReadWritePaths=/home/unitree/unitree_sdk2-main
ReadWritePaths=/home/unitree/unitree_sdk2
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target 