# control_robot.service
[Unit]
Description=Control Robot Service
After=network.target docker.service
Wants=docker.service
Requires=network.target

[Service]
Type=simple
User=unitree
Group=unitree
WorkingDirectory=/home/unitree/control_robot
ExecStart=/home/unitree/control_robot/start_h1_unified.sh
ExecStop=/home/unitree/control_robot/stop_h1.sh
Restart=always
StandardOutput=journal
StandardError=journal
SyslogIdentifier=control_robot

# Переменные окружения
Environment=NODE_ENV=production
Environment=HOME=/home/unitree
Environment=SYSTEMD_EXEC_PID=1

# Права доступа
ReadWritePaths=/home/unitree
ReadWritePaths=/home/unitree/control_robot
ReadWritePaths=/home/unitree/unitree_sdk2-main
ReadWritePaths=/home/unitree/unitree_sdk2

[Install]
WantedBy=multi-user.target 